# monitors/uptime/uptime-monitor.timer
[Unit]
Description=Run uptime dual-vantage monitor every 10 min

[Timer]
OnBootSec=2min
OnUnitActiveSec=10min
AccuracySec=30s
Persistent=true

[Install]
WantedBy=timers.target
